# misc/blinky-tape.service.in
[Unit]
Description=Blinky Tape

[Service]
ExecStart=@PREFIX@/bin/blinky-tape
Restart=on-failure

[Install]
WantedBy=multi-user.target
