# misc/blinky.service.in
[Unit]
Description=Blinky

[Service]
ExecStart=@PREFIX@/bin/blinky
Restart=on-failure

[Install]
WantedBy=multi-user.target
